# byg.service.example
[Unit]
Description=Byg skriverier med hugo
[Service]
ExecStart=<stien til byg.sh>
Type=oneshot
