# byg.service
[Unit]
Description=Byg skriverier med hugo
[Service]
ExecStart=/home/prime/services/skriverier/byg.sh
Type=oneshot
