# byg.service
[Unit]
Description=Byg skriverier med hugo
[Service]
ExecStart=/usr/local/bin/byg.sh
Type=oneshot
